# package/DEBIAN/prerm
#!/bin/sh

if [ -x "/etc/init.d/netsav" ]; then
  invoke-rc.d netsav stop || exit $?
fi

exit 0
